# root/etc/services.d/airupnp/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /usr/local/bin/airupnp -Z -x /config/airupnp.xml
